# squashfs-root/etc/hotplug.d/iface/11-customqos
#!/bin/sh
[ "$ACTION" = ifup ] && /bin/generaterule.sh interface "$INTERFACE"
