# squashfs-root/usr/lib/opkg/info/fota_git.list
/usr/bin/fota_git
/etc/init.d/fota_git.init
