# squashfs-root/usr/lib/opkg/info/heart_beat.list
/etc/rc.d/S97heart_beat
/etc/init.d/heart_beat
/usr/bin/heart_beat
