# squashfs-root/usr/lib/opkg/info/init_status.list
/etc/init.d/init_status
/etc/rc.d/S96init_status
/usr/bin/init_status
