# squashfs-root/usr/lib/opkg/info/isc-dhcp-server-ipv4.list
/usr/sbin/generate_dhcpd_lease_file
/etc/dhcpd.conf
/etc/config/dhcpd
/usr/sbin/dhcpd
/etc/init.d/dhcpd
