# squashfs-root/usr/lib/opkg/info/keybd_hard_reset.list
/etc/rc.d/S99keybd_hard_reset
/etc/init.d/keybd_hard_reset
/usr/bin/keybd_hard_reset
