# squashfs-root/usr/lib/opkg/info/radius_h_b.list
/etc/rc.d/radius_h_b
/etc/init.d/radius_h_b
/usr/bin/radius_h_b
/etc/init.d/radius_coa
/usr/bin/radius_coa
