# squashfs-root/usr/lib/opkg/info/rrdtool1.list
/usr/bin/rrdtool
/usr/bin/rrdupdate
/etc/rrd/graphcpu
/etc/init.d/rrd
/etc/rrd/graphinterface
/etc/rrd/loadrrd
