# squashfs-root/usr/lib/opkg/info/wifilith.list
/etc/wifilith/www/comfast/css/plugmenu2.png
/etc/wifilith/www/js/qrcode.js
/etc/wifilith/www/img/portal_pic_config.json
/etc/wifilith/www/js/respond.min.js
/etc/wifilith/www/img/2.jpg
/etc/wifilith/www/fonts/iconfont.woff
/etc/wifilith/www/traffic.html
/etc/wifilith/www/comfast/css/plugmenu3.png
/etc/wifilith/www/comfast/js/content.js
/etc/wifilith/www/no-internet-connection.html
/etc/wifilith/www/wxpay.html
/etc/wifilith/www/fail.html
/etc/wifilith/www/mdevice_scan.html
/etc/wifilith/www/comfast/clear.html
/etc/wifilith/www/img/3.jpg
/etc/wifilith/www/comfast/js/radius.js
/etc/wifilith/www/fonts/iconfont.eot
/etc/wifilith/www/weixin.html
/etc/wxdirect/black_useragent
/etc/wifilith/www/comfast/css/plugmenu4.png
/etc/wifilith/www/comfast/js/route/bhw.js
/etc/wifilith/www/comfast/js/route/ros.js
/etc/wifilith/www/favicon.ico
/etc/wifilith/www/img/4.jpg
/etc/wifilith/www/fonts/iconfont.ttf
/etc/wifilith/www/js/traffic.js
/etc/wifilith/www/sms.html
/etc/wifilith/www/fonts/iconfont.svg
/etc/wifilith/www/comfast/js/route/wayos.js
/etc/init.d/wifilith
/etc/wifilith/www/comfast/css/logo.png
/etc/wifilith/www/passwd.html
/etc/wifilith/whitelist.ip
/etc/wifilith/www/css/sms.css
/usr/bin/wifilith
/etc/wifilith/www/img/5.jpg
/etc/nginx/ssl/comfastlogin.key
/etc/wifilith/www/comfast/css/main.css
/etc/wifilith/www/comfast/fail.html
/etc/wifilith/www/img/weinxin-pay.png
/etc/config/wifilith
/etc/wifilith/www/js/bootstrap.min.js
/etc/wifilith/www/js/html5shiv.min.js
/etc/wifilith/www/comfast/js/route/panabit.js
/usr/bin/lithquery
/etc/wifilith/www/img/6.jpg
/etc/wifilith/www/js/local.js
/etc/wifilith/www/css/wxpay.css
/etc/wifilith/www/mobile_page.html
/etc/wifilith/www/comfast/css/passwd_input_1.png
/etc/wifilith/www/css/style.css
/etc/wifilith/www/account.html
/etc/wifilith/www/js/wechat.js
/etc/wifilith/www/comfast/js/route/global_login.js
/etc/wifilith/www/comfast/js/global.js
/etc/wifilith/www/mdevice_portal.html
/etc/wifilith/www/comfast/success.html
/etc/wifilith/www/img/7.jpg
/etc/wifilith/www/mdevice_fans.html
/etc/wifilith/www/css/wechat.css
/etc/wifilith/whitelist.domain
/etc/wifilith/www/comfast/js/route/huawei.js
/etc/rc.d/S96wifilith
/etc/wifilith/www/comfast/server.html
/etc/wifilith/www/comfast/css/passwd_input_2.png
/etc/wifilith/www/js/sms.js
/etc/wifilith/www/css/local.css
/etc/hotplug.d/iface/80-wifilith
/etc/wifilith/www/css/passwd.css
/etc/nginx/ssl/comfastlogin.pem
/etc/wifilith/www/comfast/js/route/ikuai.js
/usr/bin/wifilith-preremove
/etc/wifilith/www/comfast/js/jquery_1.7.2.js
/etc/wifilith/www/img/loading_img.gif
/etc/wifilith/www/mdevice_fast.html
/etc/wifilith/www/img/8.jpg
/etc/wifilith/www/js/jquery.flexslider-min.js
/etc/wifilith/www/css/bootstrap.css
/etc/wifilith/www/pc_page.html
/etc/wifilith/www/comfast/index.html
/etc/wifilith/www/js/wechatpay.js
/etc/wifilith/www/js/jquery.touchSwipe.min.js
/etc/wifilith/www/comfast/mianze.html
/etc/wifilith/www/css/traffic.css
/etc/wifilith/www/comfast/weblogin.html
/etc/wifilith/www/comfast/css/plugmenu1.png
/usr/bin/wifilith-postinstall
/etc/wifilith/www/local.html
/etc/wifilith/www/js/jquery-1.11.1.min.js
/etc/wifilith/www/css/account.css
/etc/wifilith/www/js/account.js
/etc/wifilith/www/img/1.jpg
/etc/wifilith/www/js/passwd.js
